# projects/telegram-automation-bot/packaging/macOS_LaunchAgent.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.telegram-automation-bot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/absolute/path/to/telegram-automation-bot/scripts/run.command</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/absolute/path/to/telegram-automation-bot</string>
  <key>StandardOutPath</key>
  <string>/tmp/telegram-automation-bot.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/telegram-automation-bot.err.log</string>
</dict>
</plist>

